# systemd/rewrite-bridge.service
[Unit]
Description=HK Ollama Rewrite Bridge (Express)
After=network.target

[Service]
Type=simple
User=hsadmin
WorkingDirectory=/workspace/hk-ollama-rewrite-bridge
ExecStart=/usr/bin/node /workspace/hk-ollama-rewrite-bridge/server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production
# Optional runtime tuning file (create /etc/default/rewrite-bridge)
# EnvironmentFile=/etc/default/rewrite-bridge

[Install]
WantedBy=multi-user.target
